# overlay.d/95k3s/usr/lib/systemd/system/install-k3s.service.d/11-skip-selinux-rpm.conf
[Service]
Environment=INSTALL_K3S_SKIP_SELINUX_RPM=true